# esphome_proxy/rootfs/etc/s6-overlay/s6-rc.d/nginx/finish
#!/command/with-contenv bashio
# ==============================================================================
# Handle nginx exit
# ==============================================================================
bashio::log.warning "NGINX has exited with code ${1}"
