# esphome_proxy/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
bashio::log.info "Testing NGINX config..."
nginx -t -c /etc/nginx/nginx.conf
bashio::log.info "Starting NGINX..."
exec nginx -c /etc/nginx/nginx.conf