# esphome_proxy/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Start nginx
# ==============================================================================
bashio::log.info "Starting NGINX..."
exec nginx -c /etc/nginx/nginx.conf
