# esphome_proxy/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Generate nginx ingress config dynamically from add-on options + supervisor API
# ==============================================================================

bashio::log.info "Configuring nginx for ESPHome proxy..."

# Read user-configured ESPHome URL from add-on options
ESPHOME_URL=$(bashio::config 'esphome_url')

# Strip trailing slash if present
ESPHOME_URL="${ESPHOME_URL%/}"

# Get the ingress path assigned by the Supervisor
INGRESS_PATH=$(bashio::addon 'ingress_path')

bashio::log.info "ESPHome URL: ${ESPHOME_URL}"
bashio::log.info "Ingress path: ${INGRESS_PATH}"

# Write the nginx server block for ingress
cat > /etc/nginx/servers/ingress.conf << NGINX
server {
    listen 5000 default_server;

    include /etc/nginx/includes/server_params.conf;

    location / {
        # Only allow requests from the HA Supervisor/Core
        allow 172.30.32.2;
        allow 172.30.32.1;
        deny all;

        proxy_pass ${ESPHOME_URL};

        # Tell ESPHome what base path it is being served from
        proxy_set_header X-Ingress-Path ${INGRESS_PATH};
        proxy_set_header Host \$http_host;

        include /etc/nginx/includes/proxy_params.conf;
    }
}
NGINX

bashio::log.info "nginx ingress config written."
