# esphome_proxy/rootfs/etc/services.d/nginx/finish
#!/command/with-contenv bashio
bashio::log.warning "NGINX exited with code ${1}"