# esphome_proxy/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting NGINX..."
nginx -c /etc/nginx/nginx.conf