# esphome_proxy/rootfs/etc/services.d/nginx/run
#!/bin/sh
echo "[$(date '+%H:%M:%S')] INFO: Starting NGINX..."
exec /usr/sbin/nginx -c /etc/nginx/nginx.conf